# bot/config/supervisord.conf
[supervisord]
nodaemon=true                    ; Run in foreground (required for Docker PID 1)
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; ── Process 1: Virtual Display (Xvfb) ─────────────────────────────────────────
; Must start FIRST — Chrome needs the display ready
[program:xvfb]
command=Xvfb :99 -screen 0 1280x720x24 -ac +extension GLX +render -noreset
autorestart=true
priority=100                     ; Lower number = starts first
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb-err.log
environment=HOME="/root"

; ── Process 2: PulseAudio (Virtual Audio) ─────────────────────────────────────
; Starts after Xvfb. Sets up virtual speaker + mic.
[program:pulseaudio]
command=pulseaudio
    --exit-idle-time=-1
    --log-level=warn
    --log-target=stderr
autorestart=true
priority=200                     ; Starts second
stdout_logfile=/var/log/supervisor/pulseaudio.log
stderr_logfile=/var/log/supervisor/pulseaudio-err.log
environment=HOME="/root",DISPLAY=":99"
startsecs=2                      ; Wait 2s before marking as "started"

; ── Process 3: Bot Script ─────────────────────────────────────────────────────
; Starts last, after display and audio are ready
[program:bot]
command=python3 /app/src/bot.py
autorestart=false                ; Do NOT restart — if bot exits, meeting is done
priority=300                     ; Starts last
stdout_logfile=/var/log/supervisor/bot.log
stderr_logfile=/var/log/supervisor/bot-err.log
environment=
    HOME="/root",
    DISPLAY=":99",
    PULSE_SERVER="unix:/tmp/pulse/native"
startsecs=5                      ; Give Xvfb + PulseAudio 5s to fully start
